# contrib/upstart/tutum-agent.conf
#!upstart
start on runlevel [2345]
stop on runlevel [06]

script
	exec /usr/bin/tutum-agent
end script
